# config/crontab
# Run 13 minutes after the hour.
* * * * * /Users/jepollock/dev/ai_masters/config/kopia_sync >/tmp/kopia_sync.log 2>&1
