# config/crontab
# Run 13 minutes after the hour.
13 * * * * /Users/jepollock/dev/ai_masters/config/kopia_sync
